# os/systemd/arcade-input.service
[Unit]
Description=Arcade Input Service
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=root
WorkingDirectory=/opt/arcade/service
EnvironmentFile=-/opt/arcade/os/.env.arcade-service

Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/arcade1/.Xauthority
Environment=GAME_PACKAGE_KEY_HEX=e93e36263af461dd686403aa98b20a1e084818359af81258fdb9b17828940620
Environment=ARCADE_RUNTIME_GAMES_DIR=/dev/shm/arcade-games
Environment=ARCADE_RUNTIME_DIR=/opt/arcade
Environment=ARCADE_UI_DIST_DIR=/opt/arcade/ui/dist
Environment=UINPUT_HELPER_PATH=/opt/arcade/bin/uinput-helper
Environment=RETROARCH_TTY_MODE=1
Environment=ARCADE_UI_VT=2
Environment=ARCADE_GAME_VT=1
Environment=RETROARCH_PRIMARY_INPUT=P1
Environment=RETROARCH_CONFIG_PATH=/opt/arcade/os/retroarch.cfg
Environment=CASINO_MENU_EXITS_RETROARCH=1
Environment=ARCADE_RETRO_OSD_COOLDOWN_MS=250
Environment=ARCADE_RETRO_OSD_PROMPT_INTERVAL_MS=400

ExecStartPre=/opt/arcade/os/bin/arcade-audio-init.sh
ExecStart=/usr/bin/node /opt/arcade/service/input.bundle.cjs
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
